# .flazhost/supervisord.conf
; FlazHost Runner — supervisor: web + queue worker + scheduler.
; CATATAN multi-replica: kalau app di-scale > 1 replica, [scheduler] jalan di
; tiap replica. Pakai ->onOneServer() pada scheduled task (butuh cache
; Redis/database) supaya cron tidak tereksekusi dobel.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:apache]
command=apache2-foreground
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
command=php /var/www/html/artisan schedule:work
user=www-data
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue-worker]
command=php /var/www/html/artisan queue:work --tries=3 --max-time=3600 --sleep=3
user=www-data
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
